# caserver-configure-iptables.service
[Unit]
Description=configure iptables for multiple CA servers

After=network-online.target

[Install]
WantedBy=multi-user.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/epics-configure-iptables start
ExecStop=/usr/local/bin/epics-configure-iptables stop
RemainAfterExit=yes
StandardError=journal
